# pkg/debian/postinst
#!/bin/sh
set -e

RPROFILE="/etc/R/Rprofile.site"
MARKER="# rapt: enable system package integration"

case "$1" in
    configure)
        # Enable systemd socket
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable raptd.socket
            systemctl start raptd.socket
        fi

        # Add rapt::enable() to Rprofile.site
        if [ -f "$RPROFILE" ]; then
            if ! grep -q "$MARKER" "$RPROFILE"; then
                cat >> "$RPROFILE" << 'EOF'

# rapt: enable system package integration
if (requireNamespace("rapt", quietly = TRUE)) rapt::enable()
EOF
            fi
        else
            # Create Rprofile.site if it doesn't exist
            mkdir -p "$(dirname "$RPROFILE")"
            cat > "$RPROFILE" << 'EOF'
# rapt: enable system package integration
if (requireNamespace("rapt", quietly = TRUE)) rapt::enable()
EOF
        fi
        ;;
esac

#DEBHELPER#

exit 0
