# pkg/debian/rapt.raptd.socket
[Unit]
Description=rapt daemon socket

[Socket]
ListenStream=/run/raptd.sock
SocketMode=0666

[Install]
WantedBy=sockets.target
