# daemon/raptd.service
[Unit]
Description=rapt daemon - R package installation via apt
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/raptd
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
