# debian/postrm
#!/bin/sh
set -e

RPROFILE="/etc/R/Rprofile.site"

case "$1" in
    remove|purge)
        # Stop and disable systemd socket
        if [ -d /run/systemd/system ]; then
            systemctl stop raptd.socket 2>/dev/null || true
            systemctl stop raptd.service 2>/dev/null || true
            systemctl disable raptd.socket 2>/dev/null || true
        fi

        # Remove rapt lines from Rprofile.site
        if [ -f "$RPROFILE" ]; then
            # Remove the rapt block (marker line and following line)
            sed -i '/^# rapt: enable system package integration$/,/^if (requireNamespace("rapt".*$/d' "$RPROFILE"
            # Clean up any trailing blank lines
            sed -i -e :a -e '/^\n*$/{$d;N;ba' -e '}' "$RPROFILE" 2>/dev/null || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
