# debian/rapt/DEBIAN/postinst
#!/bin/sh
set -e

RPROFILE="/etc/R/Rprofile.site"
MARKER="# rapt: enable system package integration"

case "$1" in
    configure)
        # Enable systemd socket
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable raptd.socket
            systemctl start raptd.socket
        fi

        # Add rapt::enable() to Rprofile.site
        if [ -f "$RPROFILE" ]; then
            if ! grep -q "$MARKER" "$RPROFILE"; then
                cat >> "$RPROFILE" << 'EOF'

# rapt: enable system package integration
if (requireNamespace("rapt", quietly = TRUE)) rapt::enable()
EOF
            fi
        else
            # Create Rprofile.site if it doesn't exist
            mkdir -p "$(dirname "$RPROFILE")"
            cat > "$RPROFILE" << 'EOF'
# rapt: enable system package integration
if (requireNamespace("rapt", quietly = TRUE)) rapt::enable()
EOF
        fi
        ;;
esac

# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'raptd.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'raptd.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'raptd.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'raptd.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'raptd.socket' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'raptd.socket'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'raptd.socket' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'raptd.socket' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'raptd.service' 'raptd.socket' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
